# timer.conf
#!upstart
description "Timer App"

start on started mountall
stop on shutdown

# Automatically Respawn:
respawn
respawn limit 99 5

env NODE_ENV=production

exec /home/bitnami/timer/app/runl.sh
